# plex/rootfs/etc/services.d/plex/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Plex Media Server
# Runs the Plex Media Server
# ==============================================================================
if bashio::config.has_value 'mounts'; then
  while read -r mountspec; do
    readarray -d "|" -t m <<< "${mountspec}"
    remote=`echo ${m[1]} | xargs`
    local=`echo ${m[0]} | xargs`

    bashio::log.info "Mounting \"${remote}\" to \"/mnt/${local}\""

    mkdir -p "/mnt/${local}"
    mount -v "${remote}" "/mnt/${local}"
  done <<< "$(bashio::config 'mounts')"
fi


bashio::log.info 'Starting the Plex Media Server...'

export LD_LIBRARY_PATH=/usr/lib/plexmediaserver
export PLEX_MEDIA_SERVER_HOME=/usr/lib/plexmediaserver
export PLEX_MEDIA_SERVER_MAX_PLUGIN_PROCS=6
export PLEX_MEDIA_SERVER_INFO_DEVICE="Hass.io"
export PLEX_MEDIA_SERVER_APPLICATION_SUPPORT_DIR=/data

# Run the Plex Media Server
exec /usr/lib/plexmediaserver/Plex\ Media\ Server
